# telegram-bot/contract-bot.service
[Unit]
Description=QiAll Design Contract Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/contract-bot
ExecStart=/opt/contract-bot/venv/bin/python bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
